# deploy/ticketsystem.service
[Unit]
Description=IT-Ticketsystem (Gunicorn)
After=network.target

[Service]
# WICHTIG: User/Pfade ggf. anpassen
User=ticket
Group=www-data
WorkingDirectory=/home/ticket/ticket-system
EnvironmentFile=/home/ticket/ticket-system/.env
ExecStart=/home/ticket/ticket-system/.venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/home/ticket/ticket-system/gunicorn.sock \
    --access-logfile - \
    --error-logfile - \
    ticketsystem.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
